# systemd/claude-line-notifier.service
[Unit]
Description=Claude-Code LINE Notifier
After=network.target claude-a2a.service claude-webhook.service
PartOf=claude-a2a.service

[Service]
Type=simple
User=planj
WorkingDirectory=/home/planj/Claude-Code-Communication
Environment="PATH=/home/planj/Claude-Code-Communication/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
EnvironmentFile=/home/planj/Claude-Code-Communication/.env
ExecStart=/home/planj/Claude-Code-Communication/venv/bin/python3 /home/planj/Claude-Code-Communication/integrations/line_notifier.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
